# raw/domains/codex-native-memory-governance/repos/codex-memory-kit/launchd/com.codex.mult-agent.feishu-bridge.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.codex.mult-agent.feishu-bridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/&lt;user&gt;/project/&lt;legacy-upstream-root&gt;/bin/feishu-codex-bridge-local.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/&lt;user&gt;/project/&lt;legacy-upstream-root&gt;</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/&lt;user&gt;</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>CODEX_FEISHU_BRIDGE_ENV</key>
    <string>/Users/&lt;user&gt;/.config/&lt;legacy-upstream-root&gt;/feishu-bridge.env</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/&lt;user&gt;/project/&lt;legacy-upstream-root&gt;/.omx/logs/feishu-codex-bridge.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/&lt;user&gt;/project/&lt;legacy-upstream-root&gt;/.omx/logs/feishu-codex-bridge.err.log</string>
</dict>
</plist>
